# package/moneta-backup.service
[Unit]
Description=Backup Moneta

[Service]
Type=oneshot
ExecStart=/bin/bash /var/backups/backup-moneta.sh
